# crontab
# Run every 30 minutes
*/30 * * * * . /env.sh; /usr/local/bin/python /app/main.py >> /var/log/cron.log 2>&1
